# alpha-panel/supervisor/terminal.conf
[program:terminal-serve]
command=/usr/local/bin/php /var/www/AlphaPanel/httpdocs/artisan terminal:serve --port=2999
autostart=true
autorestart=true